# src/dist/omero-ms-thumbnail.service
[Unit]
Description=OMERO thumbnail microservice server
Requires=multi-user.target
After=multi-user.target
After=network.service

[Service]
Type=simple
Environment="JAVA_OPTS=-Dlogback.configurationFile=/opt/omero/OMERO.ms/omero-ms-thumbnail/current/conf/logback.xml"
WorkingDirectory=/opt/omero/OMERO.ms/omero-ms-thumbnail/current
ExecStart=/opt/omero/OMERO.ms/omero-ms-thumbnail/current/bin/omero-ms-thumbnail
User=omero
Group=omero
Restart=no
LimitNPROC=16384
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
